# docker/unified/supervisord.conf
; Supervisor manages all the per-language runtimes inside the unified
; container. Each program writes its log to /var/log/cis/<name>.log and
; emits to the container's stdout/stderr so Railway captures them.

[supervisord]
nodaemon=true
user=root
logfile=/var/log/cis/supervisord.log
pidfile=/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=/usr/sbin/php-fpm8.2 --nodaemonize --fpm-config /etc/php/8.2/fpm/php-fpm.conf
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bcgov-api]
command=/opt/venv/bin/python -m uvicorn app:app --host 127.0.0.1 --port 8001 --log-level info
directory=/opt/api
priority=30
autorestart=true
environment=PYTHONUNBUFFERED=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:auth-bridge]
command=/usr/bin/npx tsx auth-bridge/server.ts
directory=/opt/integration
priority=40
autorestart=true
environment=AUTH_BRIDGE_PORT=8090
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:metabase-sso]
command=/usr/bin/npx tsx metabase-sso/server.ts
directory=/opt/integration
priority=50
autorestart=true
environment=METABASE_SSO_PORT=8091
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
